# init/dnska.service
[Unit]
Description=Dnska
Wants=network.target

[Service]
Type=simple
ExecStart=/app/dnska

[Install]
WantedBy=multi-user.target
